# deploy/systemd/minerva-mint.service
[Unit]
Description=Minerva Mint — Ark-backed Cashu server
Documentation=https://github.com/Z0rlord/minervamnt
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=<USER>
Group=<USER>
WorkingDirectory=/opt/minervamnt
EnvironmentFile=/opt/minervamnt/.env
ExecStart=/opt/minervamnt/target/release/minerva-mint
Restart=on-failure
RestartSec=5

# Hardening (adjust ReadWritePaths for your data dir)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
